# images/unit3d/rootfs/services/rsyslogd/run
#!/bin/sh

mkdir -p /var/spool/rsyslog
rm -f /var/run/rsyslogd.pid
exec rsyslogd -n -f /etc/rsyslog/rsyslog.conf
