# images/unit3d/rootfs/services/rsyslogd/run
#!/bin/sh

mkdir -p /var/spool/rsyslog

exec rsyslogd -n -f /etc/rsyslog/rsyslog.conf
